# supervisord/supervisord_housekeeping.conf
[program:prizm_housekeeping]
command=python /home/pi/daq_2018/prizm_housekeeping.py -o "antenna","res50","short","noise","res100","open" -a 23,1 -r 24,1 -s 25,1 -n 5,1 -R 26,1 -O 12,1 --reset 18 -m 21 -t 1 -c /home/pi/daq_2018/config_tempsensor.txt
process_name='prizm_housekeeping'
startsecs=5
stdout_logfile=/home/pi/supervisord_logs/prizm_housekeeping_stdout.log
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=10
stderr_logfile=/home/pi/supervisord_logs/prizm_housekeeping_stderr.log
stderr_logfile_maxbytes=1MB
stderr_logfile_backups=10

[eventlistener:prizm_daq_listener]
command=/home/pi/daq_2018/supervisord/prizm_daq_listener.py
process_name='prizm_daq_listener'
events=PROCESS_STATE
stdout_logfile=/home/pi/supervisord_logs/prizm_daq_listener_stdout.log
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=10
stderr_logfile=/home/pi/supervisord_logs/prizm_daq_listener_stderr.log
stderr_logfile_maxbytes=1MB
stderr_logfile_backups=10

[supervisord]
logfile=/home/pi/supervisord_logs/supervisord.log
logfile_maxbytes=1MB
logfile_backups=10
loglevel=info
pidfile=/home/pi/supervisord_logs/supervisord.pid
